# todobg.service
[Service]
ExecStart=/usr/bin/todobg
Restart=always
