# todobg-recreate.service
[Service]
Type=exec
ExecStart=/usr/bin/todobg --recreate
